# .ci/docker-ansimulator-template/system/profile.d/custom.sh
# keep ls classic behavior
export QUOTING_STYLE=literal

# force ll alias to ls -a"
alias ll='ls -alF'
